# examples/bat.timer
[Unit]
Description=Simple battery watcher script to hibernate the laptop

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
